# backend/aws/vnc-login/supervisor-display.conf
[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 1920x1080x24 -ac +extension GLX +render -noreset
autorestart=true
priority=100
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.err

[program:openbox]
command=/usr/bin/openbox
autorestart=true
priority=150
environment=DISPLAY=":99"
stdout_logfile=/var/log/openbox.log
stderr_logfile=/var/log/openbox.err

[program:x11vnc]
# Keep x11vnc input-friendly for noVNC:
# -rfbversion 3.8: better client compatibility
# -skip_dups/-nomodtweak/-norepeat: avoid missed keystrokes and stuck modifiers
# -xkb/-repeat: ensure keyboard map + repeats are enabled
# -noxdamage: smoother updates on Xvfb
command=/usr/local/bin/x11vnc -forever -display :99 -rfbport 5900 -rfbversion 3.8 -nopw -shared -xkb -repeat -skip_dups -nomodtweak -norepeat -noxdamage -wait 10 -defer 10 -cursor most
autorestart=true
priority=200
environment=DISPLAY=":99"
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.err

[program:novnc]
command=/opt/novnc/utils/websockify/run --web /opt/novnc 80 localhost:5900
autorestart=true
priority=300
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.err

